# debian/imx6-config.bootsplash.init
#! /bin/sh
### BEGIN INIT INFO
# Provides:          bootsplash
# Required-Start:    $local_fs
# Required-Stop:     
# Default-Start:     S
# Default-Stop:      S
# Short-Description: Show this to console during boot
### END INIT INFO


PATH=/sbin:/usr/sbin:/bin:/usr/bin

echo "Booting Debian for Imx6 (SD-CARD)..." >> /dev/tty0
fbi -once -timeout 5 -noverbose -vt 6 /etc/bootsplash.png
